# ops/hetzner/systemd/cortex-gateway.service
[Unit]
Description=Cortex Gateway for Eudaemon Alpha Phase 6
After=network-online.target cortex-icp-network.service
Wants=network-online.target cortex-icp-network.service

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__DEPLOY_ROOT__/repo/cortex
EnvironmentFile=__DEPLOY_ROOT__/config/eudaemon-alpha.env
Environment=WORK_ROUTER_LOG_ROOT=__DEPLOY_ROOT__/logs/work_router
ExecStart=__DEPLOY_ROOT__/repo/cortex/target/release/cortex-gateway
Restart=on-failure
RestartSec=5
StandardOutput=append:__DEPLOY_ROOT__/logs/cortex-gateway.log
StandardError=append:__DEPLOY_ROOT__/logs/cortex-gateway.log

[Install]
WantedBy=multi-user.target
