# ops/hetzner/systemd/cortex-worker.service
[Unit]
Description=Nostra Cortex Native Worker
After=network.target cortex-icp-network.service
Wants=cortex-icp-network.service

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__DEPLOY_ROOT__/repo/nostra/worker
EnvironmentFile=__DEPLOY_ROOT__/config/eudaemon-alpha.env
ExecStart=__DEPLOY_ROOT__/repo/nostra/worker/target/release/cortex_worker
Restart=always
RestartSec=5
StandardOutput=append:__DEPLOY_ROOT__/logs/cortex-worker.log
StandardError=append:__DEPLOY_ROOT__/logs/cortex-worker.log

[Install]
WantedBy=multi-user.target
